# clusivai-ngrok.service
[Unit]
Description=Clusivai ngrok Tunnel
After=network.target clusivai-webapp.service
Requires=clusivai-webapp.service

[Service]
Type=simple
User=clusiv
WorkingDirectory=/home/clusiv/bot-recordatorios
ExecStart=/usr/bin/env bash /home/clusiv/bot-recordatorios/start-ngrok.sh
Restart=always
RestartSec=5
EnvironmentFile=/home/clusiv/bot-recordatorios/.env

[Install]
WantedBy=multi-user.target