# clusivai-webapp.service
[Unit]
Description=Clusivai Web App Server
After=network.target

[Service]
Type=simple
User=clusiv
WorkingDirectory=/home/clusiv/bot-recordatorios/bot-clusivai
ExecStart=/home/clusiv/bot-recordatorios/bot-clusivai/venv/bin/python3 /home/clusiv/bot-recordatorios/bot-clusivai/server.py
Restart=always
RestartSec=10
EnvironmentFile=/home/clusiv/bot-recordatorios/bot-clusivai/.env

[Install]
WantedBy=multi-user.target
